# systemd/btts.target
[Unit]
Description="Bluetooth test suite"
After=ofono.service
Requires=ofono.service

[Install]
WantedBy=multi-user.target
