# systemd/btts-bluez-agent.service
[Unit]
Description="Bluetooth test suite - bluez agent"
After=bluetooth.service

[Service]
Type=dbus
BusName=org.merproject.btts.BluezAgent
ExecStart=/usr/bin/btts bluez-agent --mode=server

[Install]
WantedBy=btts.target
