# systemd/btts.target
[Unit]
Description="Bluetooth test suite"
